# chameleon.timer
[Unit]
Description=Run Theme Switcher (every half hour)

[Timer]
OnCalendar=*:00,30:00
Persistent=true

[Install]
WantedBy=timers.target
