# chameleon.service
[Unit]
Description=Colorscheme Switcher (Day/Night)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/chameleon
